# alpine/etc/init.d/trmnl_srv
#!/sbin/openrc-run

name="trmnl_srv"
description="TRMNL Server"
command="/home/donmcc/trmnl_srv/serve_trmnl"
command_user="root"
pidfile="/run/${RC_SVCNAME}.pid"
command_background="yes"
output_log="/var/log/${RC_SVCNAME}.log"
error_log="/var/log/${RC_SVCNAME}.err"

depend() {
    need net
    after firewall
}

start() {
    ebegin "Starting ${name}"
    start-stop-daemon --start \
        --exec ${command} \
        --user ${command_user} \
        --make-pidfile --pidfile ${pidfile} \
        --background \
        --stdout ${output_log} \
        --stderr ${error_log}
    eend $?
}

stop() {
    ebegin "Stopping ${name}"
    start-stop-daemon --stop \
        --pidfile ${pidfile} \
        --retry TERM/30/KILL/5
    eend $?
}
